# Etisalat_USB_Modem/bin/9-cdrom.rules
#
# udev rules file 
#
#

#SUBSYSTEM!="usb_device", GOTO="phone_rules_end"
#ACTION!="add", GOTO="phone_rules_end"KERNEL=="sr*",SUBSYSTEMS=="scsi"
#ACTION!="remove", GOTO="phone_rules_end",KERNEL=="scd*",NAME="%k", ACTION=="add",

#################
#################
# disable NetworkManager/modem-manager
ATTRS{idVendor}=="19d2", ENV{ID_MM_DEVICE_IGNORE}="1"
#################
#################


#################
#################
#autorun#
#SUBSYSTEMS=="usb", ACTION=="add", SYSFS{idVendor}=="19d2",SYSFS{idProduct}=="2000", MODE="0666", RUN+="/bin/launch-gui.sh"
#SUBSYSTEMS=="usb", ACTION=="add", SYSFS{idVendor}=="19d2",SYSFS{idProduct}=="0053", MODE="0666", RUN+="/bin/launch-gui.sh"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="19d2",ATTRS{idProduct}=="2000", MODE="0666", RUN+="/bin/launch-gui.sh"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="19d2",ATTRS{idProduct}=="0154", MODE="0666", RUN+="/bin/launch-gui.sh"
SUBSYSTEMS=="usb", ACTION=="remove", RUN+="/bin/launch-gui.sh"
#################
#################


#################
#################
#guolv driver
SUBSYSTEM=="block", \
SYSFS{idVendor}=="19d2", \
SYSFS{idProduct}=="0154", \
SYMLINK+="ZTECdrom",\
RUN+="/usr/bin/eject /dev/ZTECdrom"

SUBSYSTEM=="block", \
SYSFS{idVendor}=="19d2", \
SYSFS{idProduct}=="2000", \
SYMLINK+="ZTECdrom",\
RUN+="/usr/bin/eject /dev/ZTECdrom"
#################
#################


#################
#################
#create MODEMCONNECTION
SUBSYSTEMS=="usb",KERNEL=="sr[0-9]*", GROUP="cdrom",ATTRS{idVendor}=="19d2", ATTRS{idProduct}=="2000",SYMLINK+="MODEMCONNECTION",MODE="0666", RUN+=""
SUBSYSTEMS=="usb",KERNEL=="sr[0-9]*", GROUP="cdrom",ATTRS{idVendor}=="19d2", ATTRS{idProduct}=="0154",SYMLINK+="MODEMCONNECTION",MODE="0666", RUN+=""
#################
#################

LABEL="phone_rules_end"


